# solar-control-dev/rootfs/etc/services.d/solar-control/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the solar optimiser service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

bashio::log.info "Starting solar optimiser service..."

# Check if Python is available
if ! command -v python3 &> /dev/null; then
    bashio::log.error "Python3 is not installed!"
    exit 1
fi

# Check if the program exists
if [ ! -f /usr/bin/my_program.py ]; then
    bashio::log.error "my_program.py not found at /usr/bin/my_program.py"
    exit 1
fi

# Check if the program is executable
if [ ! -x /usr/bin/my_program.py ]; then
    bashio::log.error "my_program.py is not executable"
    exit 1
fi

# Check if templates directory exists
if [ ! -d /usr/bin/templates ]; then
    bashio::log.error "templates directory not found at /usr/bin/templates"
    exit 1
fi

# Debug: List contents of /data before operations
bashio::log.info "Contents of /data before operations:"
ls -la /data || true

# Create only the necessary nginx directories
bashio::log.info "Creating nginx directories..."
mkdir -p /run/nginx
mkdir -p /tmp/nginx

# Get Home Assistant URL from supervisor
HASS_URL=$(bashio::config 'hass_url')
if [ -z "${HASS_URL}" ]; then
    HASS_URL="http://supervisor/core"
fi

# Export required environment variables
export SUPERVISOR_TOKEN="${SUPERVISOR_TOKEN}"
export HASS_URL="${HASS_URL}"
export PORT=5000

# Start Python web server with debug logging
bashio::log.info "Starting Python web server on port ${PORT}..."
exec python3 /usr/bin/my_program.py
