# solar-control/rootfs/etc/services.d/solar-control/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the solar optimiser service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

bashio::log.info "Starting solar optimiser service..."

# Check if Python is available
if ! command -v python3 &> /dev/null; then
    bashio::log.error "Python3 is not installed!"
    exit 1
fi

# Check if Nginx is available
if ! command -v nginx &> /dev/null; then
    bashio::log.error "Nginx is not installed!"
    exit 1
fi

# Check if the program exists
if [ ! -f /usr/bin/my_program ]; then
    bashio::log.error "my_program not found at /usr/bin/my_program"
    exit 1
fi

# Check if the program is executable
if [ ! -x /usr/bin/my_program ]; then
    bashio::log.error "my_program is not executable"
    exit 1
fi

# Check if templates directory exists
if [ ! -d /usr/bin/templates ]; then
    bashio::log.error "templates directory not found at /usr/bin/templates"
    exit 1
fi

# Debug: List contents of /data before operations
bashio::log.info "Contents of /data before operations:"
ls -la /data || true

# Ensure data directory exists and has correct permissions
if [ ! -d /data ]; then
    bashio::log.info "Creating data directory..."
    mkdir -p /data
fi

# Set directory permissions (755 for directories)
bashio::log.info "Setting data directory permissions..."
chown -R root:root /data
chmod 755 /data

# Set file permissions (644 for files)
bashio::log.info "Setting file permissions..."
find /data -type f -exec chmod 644 {} \;
find /data -type f -exec chown root:root {} \;

# Debug: Show final directory state
bashio::log.info "Final data directory state:"
ls -la /data || true

# Get Home Assistant URL from supervisor
HASS_URL=$(bashio::config 'hass_url')
if [ -z "${HASS_URL}" ]; then
    HASS_URL="http://supervisor/core"
fi

# Export required environment variables
export SUPERVISOR_TOKEN="${SUPERVISOR_TOKEN}"
export HASS_URL="${HASS_URL}"

# Start Nginx
bashio::log.info "Starting Nginx..."
nginx

bashio::log.info "Starting Python web server..."
exec python3 /usr/bin/my_program
